# package/deb/iptables-loader/DEBIAN/prerm
#!/bin/sh

if [ -x "/etc/init.d/iptables" ]; then
  invoke-rc.d iptables stop || exit $?
fi

exit 0
