# package/debian/postinst
#!/bin/sh
# postinst script for iptables-loader
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
  configure)
    echo '####          README          ####'
    echo 'The first time start the service with "iptables-loader test"'
    echo ' in order to prevent self block by a bad firewall rule'
    echo '####          !!!!!!          ####'

    if [ -x "/etc/init.d/iptables" ] && [ -x "$(command -v update-rc.d)" ]; then
      echo 'Add service into boot sequence'
      update-rc.d iptables defaults >/dev/null || exit $?
      #invoke-rc.d iptables start || exit $?
    else
      echo 'No implemented service manager found the firewall' >&2
      echo '  will not be automatically launch at boot' >&2
      exit 1
    fi
  ;;

  abort-upgrade|abort-remove|abort-deconfigure)
  ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0