# package/DEBIAN/postinst
#!/bin/sh

if [ -x "/etc/init.d/iptables-loader" ]; then
  update-rc.d iptables-loader defaults >/dev/null
  invoke-rc.d iptables-loader start || exit $?
fi

exit 0
