# package/DEBIAN/prerm
#!/bin/sh

if [ -x "/etc/init.d/iptables-loader" ]; then
  invoke-rc.d iptables-loader stop || exit $?
fi

exit 0
